# run/systemd/generator.late/sysinit.target.wants/kbd.service
# Automatically generated by systemd-sysv-generator

[Unit]
SourcePath=/etc/init.d/kbd
Description=LSB: Prepare console
DefaultDependencies=no
Before=sysinit.target
After=remote-fs.target

[Service]
Type=forking
Restart=no
TimeoutSec=0
IgnoreSIGPIPE=no
KillMode=process
GuessMainPID=no
RemainAfterExit=yes
SysVStartPriority=15
ExecStart=/etc/init.d/kbd start
ExecStop=/etc/init.d/kbd stop
