# run/systemd/generator.late/sysinit.target.wants/nfs-common.service
# Automatically generated by systemd-sysv-generator

[Unit]
SourcePath=/etc/init.d/nfs-common
Description=LSB: NFS support files common to client and server
DefaultDependencies=no
Before=sysinit.target shutdown.target
After=rpcbind.target time-sync.target
Conflicts=shutdown.target

[Service]
Type=forking
Restart=no
TimeoutSec=0
IgnoreSIGPIPE=no
KillMode=process
GuessMainPID=no
RemainAfterExit=yes
SysVStartPriority=8
ExecStart=/etc/init.d/nfs-common start
ExecStop=/etc/init.d/nfs-common stop
