# run/systemd/generator.late/frontol-du.service
# Automatically generated by systemd-sysv-generator

[Unit]
SourcePath=/etc/init.d/frontol-du
Description=LSB: No description
Before=shutdown.target
After=remote-fs.target systemd-journald-dev-log.socket local-fs.target network-online.target
Wants=network-online.target
Conflicts=shutdown.target

[Service]
Type=forking
Restart=no
TimeoutSec=5min
IgnoreSIGPIPE=no
KillMode=process
GuessMainPID=no
RemainAfterExit=yes
ExecStart=/etc/init.d/frontol-du start
ExecStop=/etc/init.d/frontol-du stop
