# run/systemd/generator.late/rc.service
# Automatically generated by systemd-sysv-generator

[Unit]
SourcePath=/etc/init.d/rc
Description=(null)

[Service]
Type=forking
Restart=no
TimeoutSec=5min
IgnoreSIGPIPE=no
KillMode=process
GuessMainPID=no
RemainAfterExit=yes
ExecStart=/etc/init.d/rc start
ExecStop=/etc/init.d/rc stop
