# run/systemd/generator.late/rcS.service
# Automatically generated by systemd-sysv-generator

[Unit]
SourcePath=/etc/init.d/rcS
Description=(null)

[Service]
Type=forking
Restart=no
TimeoutSec=5min
IgnoreSIGPIPE=no
KillMode=process
GuessMainPID=no
RemainAfterExit=yes
ExecStart=/etc/init.d/rcS start
ExecStop=/etc/init.d/rcS stop
