# run/systemd/generator.late/rpcbind.service
# Automatically generated by systemd-sysv-generator

[Unit]
SourcePath=/etc/init.d/rpcbind
Description=LSB: RPC portmapper replacement
Before=shutdown.target
After=network-online.target local-fs.target
Wants=network-online.target
Conflicts=shutdown.target

[Service]
Type=forking
Restart=no
TimeoutSec=5min
IgnoreSIGPIPE=no
KillMode=process
GuessMainPID=no
RemainAfterExit=yes
ExecStart=/etc/init.d/rpcbind start
ExecStop=/etc/init.d/rpcbind stop
